# arenda-zala.service
[Unit]
Description=Arenda Zala Conference Hall Booking System
Documentation=https://github.com/username/conference-zal
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/conference-zal
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=8027
EnvironmentFile=/var/www/conference-zal/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arenda-zala

# Security settings (without readonly for working directory)
ReadWritePaths=/var/www/conference-zal

[Install]
WantedBy=multi-user.target
